# usr/share/harbour-dynclock/harbour-dynclock.timer
[Unit]
Description=Timer for DynClock.

[Timer]
OnBootSec=0
OnCalendar=*:0/5
OnUnitActiveSec=5m
Persistent=True

[Install]
WantedBy=timers.target
